# VerionControll/Backend/Steam-API-Backend/steam-api.service
[Unit]
Description=Steam API Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/steam-api-backend
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target